# recipes-core/systemd/systemd/audio.slice
[Unit]
Description=Audio Slice
Before=slices.target
