# recipes-core/systemd/systemd/core-services.slice
[Unit]
Description=Core Services Slice
Before=slices.target
