# recipes-core/systemd/systemd/dlt.slice
[Unit]
Description=Diagnostic Log and Trace Slice
Before=slices.target
