# recipes-core/systemd/systemd/media.slice
[Unit]
Description=Media Slice
Before=slices.target
