# recipes-core/systemd/systemd/phone.slice
[Unit]
Description=Phone Slice
Before=slices.target
