# recipes-core/systemd/systemd/speech.slice
[Unit]
Description=Speech Services Slice
Before=slices.target
