# recipes-core/systemd/systemd/swdl.slice
[Unit]
Description=Software Download Slice
Before=slices.target
