# recipes-core/systemd/systemd/ui-apps.slice
[Unit]
Description=UI Apps Slice
Before=slices.target
