# recipes-core/systemd/systemd/vehicle-functions.slice
[Unit]
Description=Vehicle Functions Slice
Before=slices.target
