# recipes-core/systemd/systemd/video.slice
[Unit]
Description=Video Slice
Before=slices.target
